# build/package/etc/systemd/system/preppi.service
[Unit]
Description=preppi
Requires=local-fs.target
After=local-fs.target

[Service]
ExecStart=/usr/local/bin/preppi prepare -reboot
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
