# build/package/DEBIAN/postinst
#!/bin/sh
/bin/systemctl enable preppi.service